# src/shell/crontab
* * * * * /usr/local/bin/php /var/www/html/admin/cli/cron.php >> /var/custom_logs/cron.log
